# templates/default/etc/supervisor/pgbouncer.conf.erb
[program:pgbouncer]
command=/usr/local/bin/pgbouncer /etc/pgbouncer.ini
user=postgres
autostart=true
autorestart=true
